# config/supervisor.conf
[program:doctrine]
command=/home/ben/doctrine/venv/bin/gunicorn -b localhost:8004 -w 4 main:app
directory=/home/ben/doctrine
user=ben
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true